# raspbian_changes/cd-action.service
[Unit]
Description=Rip and/or play music

[Service]
User=pi
ExecStart=/home/pi/cd-action.sh

[Install]
WantedBy=multi-user.target
